# ubuntu_battery/docker/crontab
* * * * * (sleep 30; sh /read.sh)
* * * * * (sh /read.sh)
